# undertow.service
[Unit]
Description=The Undertow Intelligence System
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/undertow
Environment=PATH=/home/ubuntu/undertow/venv/bin
EnvironmentFile=/home/ubuntu/undertow/.env
ExecStart=/home/ubuntu/undertow/venv/bin/python -m undertow.runner
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=undertow

[Install]
WantedBy=multi-user.target

